# crawling/supervisord.conf
[program:1f602]
command=python ./crawling.py run 1f602
autostart=true
autorestart=true
stderr_logfile=./Logs/1f602/error.log
stdout_logfile=./Logs/1f602/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f612]
command=python ./crawling.py run 1f612
autostart=true
autorestart=true
stderr_logfile=./Logs/1f612/error.log
stdout_logfile=./Logs/1f612/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f629]
command=python ./crawling.py run 1f629
autostart=true
autorestart=true
stderr_logfile=./Logs/1f629/error.log
stdout_logfile=./Logs/1f629/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f62d]
command=python ./crawling.py run 1f62d
autostart=true
autorestart=true
stderr_logfile=./Logs/1f62d/error.log
stdout_logfile=./Logs/1f62d/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f60d]
command=python ./crawling.py run 1f60d
autostart=true
autorestart=true
stderr_logfile=./Logs/1f60d/error.log
stdout_logfile=./Logs/1f60d/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f614]
command=python ./crawling.py run 1f614
autostart=true
autorestart=true
stderr_logfile=./Logs/1f614/error.log
stdout_logfile=./Logs/1f614/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f44c]
command=python ./crawling.py run 1f44c
autostart=true
autorestart=true
stderr_logfile=./Logs/1f44c/error.log
stdout_logfile=./Logs/1f44c/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f60a]
command=python ./crawling.py run 1f60a
autostart=true
autorestart=true
stderr_logfile=./Logs/1f60a/error.log
stdout_logfile=./Logs/1f60a/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:2764]
command=python ./crawling.py run 2764
autostart=true
autorestart=true
stderr_logfile=./Logs/2764/error.log
stdout_logfile=./Logs/2764/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f60f]
command=python ./crawling.py run 1f60f
autostart=true
autorestart=true
stderr_logfile=./Logs/1f60f/error.log
stdout_logfile=./Logs/1f60f/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f601]
command=python ./crawling.py run 1f601
autostart=true
autorestart=true
stderr_logfile=./Logs/1f601/error.log
stdout_logfile=./Logs/1f601/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f3b6]
command=python ./crawling.py run 1f3b6
autostart=true
autorestart=true
stderr_logfile=./Logs/1f3b6/error.log
stdout_logfile=./Logs/1f3b6/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f633]
command=python ./crawling.py run 1f633
autostart=true
autorestart=true
stderr_logfile=./Logs/1f633/error.log
stdout_logfile=./Logs/1f633/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f4af]
command=python ./crawling.py run 1f4af
autostart=true
autorestart=true
stderr_logfile=./Logs/1f4af/error.log
stdout_logfile=./Logs/1f4af/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f634]
command=python ./crawling.py run 1f634
autostart=true
autorestart=true
stderr_logfile=./Logs/1f634/error.log
stdout_logfile=./Logs/1f634/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f60c]
command=python ./crawling.py run 1f60c
autostart=true
autorestart=true
stderr_logfile=./Logs/1f60c/error.log
stdout_logfile=./Logs/1f60c/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f603]
command=python ./crawling.py run 1f603
autostart=true
autorestart=true
stderr_logfile=./Logs/1f603/error.log
stdout_logfile=./Logs/1f603/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f64c]
command=python ./crawling.py run 1f64c
autostart=true
autorestart=true
stderr_logfile=./Logs/1f64c/error.log
stdout_logfile=./Logs/1f64c/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f495]
command=python ./crawling.py run 1f495
autostart=true
autorestart=true
stderr_logfile=./Logs/1f495/error.log
stdout_logfile=./Logs/1f495/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f611]
command=python ./crawling.py run 1f611
autostart=true
autorestart=true
stderr_logfile=./Logs/1f611/error.log
stdout_logfile=./Logs/1f611/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f605]
command=python ./crawling.py run 1f605
autostart=true
autorestart=true
stderr_logfile=./Logs/1f605/error.log
stdout_logfile=./Logs/1f605/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f64f]
command=python ./crawling.py run 1f64f
autostart=true
autorestart=true
stderr_logfile=./Logs/1f64f/error.log
stdout_logfile=./Logs/1f64f/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f615]
command=python ./crawling.py run 1f615
autostart=true
autorestart=true
stderr_logfile=./Logs/1f615/error.log
stdout_logfile=./Logs/1f615/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f618]
command=python ./crawling.py run 1f618
autostart=true
autorestart=true
stderr_logfile=./Logs/1f618/error.log
stdout_logfile=./Logs/1f618/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:2665]
command=python ./crawling.py run 2665
autostart=true
autorestart=true
stderr_logfile=./Logs/2665/error.log
stdout_logfile=./Logs/2665/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f610]
command=python ./crawling.py run 1f610
autostart=true
autorestart=true
stderr_logfile=./Logs/1f610/error.log
stdout_logfile=./Logs/1f610/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f481]
command=python ./crawling.py run 1f481
autostart=true
autorestart=true
stderr_logfile=./Logs/1f481/error.log
stdout_logfile=./Logs/1f481/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f61e]
command=python ./crawling.py run 1f61e
autostart=true
autorestart=true
stderr_logfile=./Logs/1f61e/error.log
stdout_logfile=./Logs/1f61e/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f648]
command=python ./crawling.py run 1f648
autostart=true
autorestart=true
stderr_logfile=./Logs/1f648/error.log
stdout_logfile=./Logs/1f648/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f62b]
command=python ./crawling.py run 1f62b
autostart=true
autorestart=true
stderr_logfile=./Logs/1f62b/error.log
stdout_logfile=./Logs/1f62b/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:270c]
command=python ./crawling.py run 270c
autostart=true
autorestart=true
stderr_logfile=./Logs/270c/error.log
stdout_logfile=./Logs/270c/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f60e]
command=python ./crawling.py run 1f60e
autostart=true
autorestart=true
stderr_logfile=./Logs/1f60e/error.log
stdout_logfile=./Logs/1f60e/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f621]
command=python ./crawling.py run 1f621
autostart=true
autorestart=true
stderr_logfile=./Logs/1f621/error.log
stdout_logfile=./Logs/1f621/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f44d]
command=python ./crawling.py run 1f44d
autostart=true
autorestart=true
stderr_logfile=./Logs/1f44d/error.log
stdout_logfile=./Logs/1f44d/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f622]
command=python ./crawling.py run 1f622
autostart=true
autorestart=true
stderr_logfile=./Logs/1f622/error.log
stdout_logfile=./Logs/1f622/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f62a]
command=python ./crawling.py run 1f62a
autostart=true
autorestart=true
stderr_logfile=./Logs/1f62a/error.log
stdout_logfile=./Logs/1f62a/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f60b]
command=python ./crawling.py run 1f60b
autostart=true
autorestart=true
stderr_logfile=./Logs/1f60b/error.log
stdout_logfile=./Logs/1f60b/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f624]
command=python ./crawling.py run 1f624
autostart=true
autorestart=true
stderr_logfile=./Logs/1f624/error.log
stdout_logfile=./Logs/1f624/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:270b]
command=python ./crawling.py run 270b
autostart=true
autorestart=true
stderr_logfile=./Logs/270b/error.log
stdout_logfile=./Logs/270b/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f637]
command=python ./crawling.py run 1f637
autostart=true
autorestart=true
stderr_logfile=./Logs/1f637/error.log
stdout_logfile=./Logs/1f637/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f44f]
command=python ./crawling.py run 1f44f
autostart=true
autorestart=true
stderr_logfile=./Logs/1f44f/error.log
stdout_logfile=./Logs/1f44f/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f440]
command=python ./crawling.py run 1f440
autostart=true
autorestart=true
stderr_logfile=./Logs/1f440/error.log
stdout_logfile=./Logs/1f440/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f52b]
command=python ./crawling.py run 1f52b
autostart=true
autorestart=true
stderr_logfile=./Logs/1f52b/error.log
stdout_logfile=./Logs/1f52b/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f623]
command=python ./crawling.py run 1f623
autostart=true
autorestart=true
stderr_logfile=./Logs/1f623/error.log
stdout_logfile=./Logs/1f623/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f608]
command=python ./crawling.py run 1f608
autostart=true
autorestart=true
stderr_logfile=./Logs/1f608/error.log
stdout_logfile=./Logs/1f608/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f613]
command=python ./crawling.py run 1f613
autostart=true
autorestart=true
stderr_logfile=./Logs/1f613/error.log
stdout_logfile=./Logs/1f613/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f494]
command=python ./crawling.py run 1f494
autostart=true
autorestart=true
stderr_logfile=./Logs/1f494/error.log
stdout_logfile=./Logs/1f494/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f498]
command=python ./crawling.py run 1f498
autostart=true
autorestart=true
stderr_logfile=./Logs/1f498/error.log
stdout_logfile=./Logs/1f498/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f3a7]
command=python ./crawling.py run 1f3a7
autostart=true
autorestart=true
stderr_logfile=./Logs/1f3a7/error.log
stdout_logfile=./Logs/1f3a7/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f64a]
command=python ./crawling.py run 1f64a
autostart=true
autorestart=true
stderr_logfile=./Logs/1f64a/error.log
stdout_logfile=./Logs/1f64a/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f609]
command=python ./crawling.py run 1f609
autostart=true
autorestart=true
stderr_logfile=./Logs/1f609/error.log
stdout_logfile=./Logs/1f609/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f480]
command=python ./crawling.py run 1f480
autostart=true
autorestart=true
stderr_logfile=./Logs/1f480/error.log
stdout_logfile=./Logs/1f480/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f616]
command=python ./crawling.py run 1f616
autostart=true
autorestart=true
stderr_logfile=./Logs/1f616/error.log
stdout_logfile=./Logs/1f616/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f604]
command=python ./crawling.py run 1f604
autostart=true
autorestart=true
stderr_logfile=./Logs/1f604/error.log
stdout_logfile=./Logs/1f604/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f61c]
command=python ./crawling.py run 1f61c
autostart=true
autorestart=true
stderr_logfile=./Logs/1f61c/error.log
stdout_logfile=./Logs/1f61c/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f620]
command=python ./crawling.py run 1f620
autostart=true
autorestart=true
stderr_logfile=./Logs/1f620/error.log
stdout_logfile=./Logs/1f620/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f645]
command=python ./crawling.py run 1f645
autostart=true
autorestart=true
stderr_logfile=./Logs/1f645/error.log
stdout_logfile=./Logs/1f645/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f4aa]
command=python ./crawling.py run 1f4aa
autostart=true
autorestart=true
stderr_logfile=./Logs/1f4aa/error.log
stdout_logfile=./Logs/1f4aa/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f44a]
command=python ./crawling.py run 1f44a
autostart=true
autorestart=true
stderr_logfile=./Logs/1f44a/error.log
stdout_logfile=./Logs/1f44a/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f49c]
command=python ./crawling.py run 1f49c
autostart=true
autorestart=true
stderr_logfile=./Logs/1f49c/error.log
stdout_logfile=./Logs/1f49c/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f496]
command=python ./crawling.py run 1f496
autostart=true
autorestart=true
stderr_logfile=./Logs/1f496/error.log
stdout_logfile=./Logs/1f496/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f499]
command=python ./crawling.py run 1f499
autostart=true
autorestart=true
stderr_logfile=./Logs/1f499/error.log
stdout_logfile=./Logs/1f499/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:1f62c]
command=python ./crawling.py run 1f62c
autostart=true
autorestart=true
stderr_logfile=./Logs/1f62c/error.log
stdout_logfile=./Logs/1f62c/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[program:2728]
command=python ./crawling.py run 2728
autostart=true
autorestart=true
stderr_logfile=./Logs/2728/error.log
stdout_logfile=./Logs/2728/output.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
numprocs=1

[supervisord]

[supervisorctl]

[inet_http_server]
port = *:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface